# debian/postinst
#!/bin/sh
# postinst script for firstboot-it4smart
#
# see: dh_installdeb(1)

set -e

case "$1" in
        configure)
        # Enable systemd-networkd-wait-online.service
        systemctl enable systemd-networkd-wait-online.service

        ;;

        abort-upgrade|abort-remove|abort-deconfigure)
            exit 0
        ;;

        *)
            echo "postinst called with unknwn argurment \'$1'" >&2
            exit 1
        ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0